# infra/packer/modules/docker/docker-consul-template.service
# Docker Configuration Watcher using Consul-Template
#
# This is a long-running daemon that continuously monitors Consul KV store for changes.
# When Consul values change, this service:
# 1. Detects the change (via consul-template's watch mechanism)
# 2. Regenerates docker-compose.yml from the .ctmpl template
# 3. Executes 'docker-compose up -d' to apply the new configuration
#
# It uses a default empty docker-compose.ctmpl template that's present on all images,
# ensuring the service works everywhere without conditionals.
# 
# Images that need to run containers override the template with their specific configuration.

[Unit]
Description=Docker Configuration Watcher (Consul-Template)
Documentation=https://github.com/hashicorp/consul-template
Requires=network-online.target consul.service docker.service
After=network-online.target consul.service docker.service

[Service]
Type=simple
User=root
Group=root

# Runs consul-template as a daemon process
ExecStart=/usr/bin/consul-template -config=/etc/consul-template/docker.consul-template.hcl

# Allows reloading the consul-template daemon itself without stopping it
# Useful when the consul-template config changes or to force re-read templates
ExecReload=/bin/kill -HUP $MAINPID

KillSignal=SIGINT
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target