# infra/packer/modules/doctl/docker-login.service
[Unit]
Description=Login to DigitalOcean Docker Registry
After=docker.service
# Only start this service in production (not during Packer builds)
# The /etc/packer_build marker exists during builds and is removed before snapshot
ConditionPathExists=!/etc/packer_build

[Service]
Type=oneshot
RemainAfterExit=yes
Environment="HOME=/root"
ExecStart=/usr/local/bin/doctl registry login

[Install]
WantedBy=multi-user.target