# infra/packer/modules/services-base/files/docker-login.service
[Unit]
Description=Login to DigitalOcean Docker Registry
After=network-online.target docker.service consul.service
Before=consul-template.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes

# Get token from Consul KV (populated during initial setup)
ExecStartPre=/bin/bash -c 'consul kv get config/digitalocean/token > /tmp/do-token 2>/dev/null || true'

# Configure doctl with token from Consul
ExecStartPre=/bin/bash -c 'if [ -s /tmp/do-token ]; then mkdir -p /root/.config/doctl && echo "access-token: $(cat /tmp/do-token)" > /root/.config/doctl/config.yaml; fi'

# Login to Docker registry
ExecStart=/bin/bash -c 'if [ -s /tmp/do-token ]; then /usr/local/bin/doctl registry login; else echo "No token in Consul KV, skipping registry login"; fi'

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target