# infra/packer/modules/alloy/alloy.service
[Unit]
Description=Grafana Alloy
Documentation=https://grafana.com/docs/alloy/
After=docker.service
Wants=docker.service

[Service]
Type=simple
User=root
# Set Loki URL - this should be configured during deployment
Environment="LOKI_URL=http://loki:3100/loki/api/v1/push"
ExecStart=/usr/bin/alloy run /etc/alloy/config.alloy
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target